# rootfs/etc/services.d/guacd/run
#!/command/with-contenv sh

echo "Starting guacamole guacd..."
s6-setuidgid root guacd -b 0.0.0.0 -f
