# rootfs/etc/services.d/postgres/run
#!/command/with-contenv bash

echo "Starting postgres..."
s6-setuidgid postgres postgres
